# config_files/supervisord-mysqld.conf
[supervidord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:mysqld]
command=/start_mysql_daemon.sh
numprocs=1
autostart=true
autorestart=true