# config_files/supervisord-mysqld.conf
[supervidord]
nodaemon=true


[program:mysqld]
command=/start_mysql_daemon.sh
numprocs=1
autostart=true
autorestart=true